# ops/systemd/nofal-panel-admin.service
[Unit]
Description=Nofal Panel Admin (FastAPI / uvicorn)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/nofal-panel
Environment=APP_MODE=admin
EnvironmentFile=-/opt/nofal-panel/.env
ExecStart=/opt/nofal-panel/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 12020 --workers 2
Restart=always
RestartSec=2
User=root
Group=root

[Install]
WantedBy=multi-user.target

